# agent-station/editorial-assistant/station.wednesday.editorial.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>station.wednesday.editorial</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/wednesdayagentic/wednesday-station/jobs/editorial-assistant/venv/bin/python</string>
    <string>/Users/wednesdayagentic/wednesday-station/jobs/editorial-assistant/server.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/wednesdayagentic/wednesday-station/jobs/editorial-assistant</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>StandardOutPath</key>
  <string>/Users/wednesdayagentic/wednesday-station/jobs/editorial-assistant/logs/server.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/wednesdayagentic/wednesday-station/jobs/editorial-assistant/logs/server.err</string>
</dict>
</plist>
