# agent-station/myhelper/station.wednesday.myhelper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>station.wednesday.myhelper</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/wednesdayagentic/wednesday-station/jobs/myhelper/venv/bin/python</string>
    <string>/Users/wednesdayagentic/wednesday-station/jobs/myhelper/bot.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/wednesdayagentic/wednesday-station/jobs/myhelper</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>StandardOutPath</key>
  <string>/Users/wednesdayagentic/wednesday-station/jobs/myhelper/logs/bot.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/wednesdayagentic/wednesday-station/jobs/myhelper/logs/bot.err</string>
</dict>
</plist>
